# src/Sh/System/ubuntu/rc.local
#!/bin/sh -e

# sudo chown root /etc/rc.local
# sudo chmod 755 /etc/rc.local
# systemctl status rc-local.service



touch /forcefsck
